# docker/php/supervisor/supervisord.conf
[program:fpm]
command=php-fpm --nodaemonize


# laravel worker example
[program:app]
command=php /var/www/html/artisan horizon
user=www-data
startsecs=0
numprocs=1
autostart=true
autorestart=true
